# root/etc/services.d/root-tool/run
#!/usr/bin/with-contenv bash

exec /opt/${COMPANY_NAME}/mediaserver/bin/root-tool-bin
